# bash/.bash/profile.d/linux.sh
#!/bin/bash
#
# 2016-02/14:jeff
#
# Platform-specific configuration for Linux hosts
#
# NOTE(jeff): This file is sourced from bash_profile

PATH="$PATH:$HOME/.config/feh/themes"
